# Automated-script/usb-photo-import@.service
[Unit]
Description=Import photos from USB device %I

[Service]
Type=oneshot
ExecStart=/home/pi/scripts/photo_import.sh %I
